# etc/s6/novnc/run
#!/bin/sh

#s6-setuidgid foo x11vnc -create -nopw
/root/noVNC/utils/launch.sh --listen 80
